# rootfs/etc/cron.d/flox
*/15 * * * * root s6-setuidgid abc /flox/backend/artisan schedule:run > /proc/1/fd/1 2>&1
