# rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# source env variables
source /etc/cont-init.d/20-setenv

exec php-fpm7.4 --nodaemonize
